# deploy/com.maxshih.opti-route.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.maxshih.opti-route.tunnel</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/cloudflared</string>
    <string>tunnel</string>
    <string>--config</string>
    <string>/Users/max_server/.cloudflared/config-opti-route.yml</string>
    <string>run</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/max_server/repo_claude/opti-route/logs/tunnel.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/max_server/repo_claude/opti-route/logs/tunnel.err.log</string>
</dict>
</plist>
